# layers/centos-7:gui:origin-1.1/@target/atomic-openshift-routewatcher.service
[Unit]
After=origin-master.service

[Service]
ExecStart=/usr/local/libexec/atomic-openshift-routewatcher.py
KillMode=process

[Install]
WantedBy=multi-user.target
