# layers/rhel-server-7:gui:ose-3.2/@target/atomic-openshift-auth.service
[Unit]
Before=atomic-openshift-master.service

[Service]
ExecStart=/usr/local/libexec/atomic-openshift-auth.py
KillMode=process

[Install]
WantedBy=multi-user.target
