# layers/rhel-server-7:gui:ose-3.0/@target/openshift-routewatcher.service
[Unit]
After=openshift-master.service

[Service]
ExecStart=/usr/local/libexec/openshift-routewatcher.py
KillMode=process

[Install]
WantedBy=multi-user.target
