# layers/rhel-server-7:gui:ose-3.1/@target/atomic-openshift-dns-intercept.service
[Unit]
After=network-online.target
Before=atomic-openshift-master.service

[Service]
ExecStart=/usr/local/libexec/atomic-openshift-dns-intercept.py
KillMode=process

[Install]
WantedBy=multi-user.target
