# layers/rhel-server-7:gui:ose-3.1/@target/atomic-openshift-master-ipcfg.service
[Unit]
After=network-online.target
Before=atomic-openshift-master.service

[Service]
Type=oneshot
ExecStart=/usr/local/libexec/atomic-openshift-master-ipcfg.py

[Install]
WantedBy=multi-user.target
