# layers/rhel-server-7:gui:rhelosp-7/@target/demobuilder-openstack-netconf-early.service
[Unit]
Before=network.target

[Service]
Type=oneshot
ExecStart=/usr/local/libexec/demobuilder-openstack-netconf-early.sh

[Install]
WantedBy=multi-user.target
