# targets/aws/@target/fix-lib-modules.service
[Unit]
Before=network.service

[Service]
Type=oneshot
ExecStart=/usr/local/libexec/fix-lib-modules.sh

[Install]
WantedBy=multi-user.target
